# usr/lib/systemd/system/start_scripts.service
[Unit]
Description=Run costum scripts on boot.

[Service]
ExecStart=/opt/init.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
