# etc/local.d/50-ariang.start
/usr/bin/aria2c --rpc-allow-origin-all -j 100 --conf-path=/etc/aria2.conf
